# examples/systemd/artigate-low.service
[Unit]
Description=Low-side ArtiGate Go module proxy/exporter
After=network-online.target
Wants=network-online.target

[Service]
User=artigate
Group=artigate
# The low side holds the signing key, so its dashboard is a privileged control
# plane. It binds to loopback by default; ArtiGate refuses to start on a
# non-loopback address unless ARTIGATE_LOW_AUTH is set (see 'artigate hashpw')
# or ARTIGATE_LOW_ALLOW_UNAUTHENTICATED=true acknowledges a trusted, TLS-
# authenticating reverse proxy in front. Drive collects from the host itself,
# an SSH tunnel, or that proxy — never expose :8080 directly.
# Environment=ARTIGATE_LOW_AUTH=admin:$argon2id$...
ExecStart=/usr/local/bin/artigate low \
  --listen 127.0.0.1:8080 \
  --root /var/lib/artigate-low \
  --export-dir /var/spool/diode-out \
  --private-key /etc/artigate/low.ed25519 \
  --upstream-goproxy https://proxy.golang.org,direct \
  --goprivate github.com/your-org/* \
  --gonosumdb github.com/your-org/*
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/var/lib/artigate-low /var/spool/diode-out /etc/artigate

[Install]
WantedBy=multi-user.target
