# examples/systemd/artigate-high.service
[Unit]
Description=High-side ArtiGate Go module repository/importer
After=network-online.target
Wants=network-online.target

[Service]
User=artigate
Group=artigate
# Serves verified, read-only mirror content. Its state-changing admin endpoints
# (POST /admin/uploads/delete, /admin/import) are restricted to loopback callers
# by default, so administer them from this host; set
# ARTIGATE_HIGH_ALLOW_REMOTE_ADMIN=on (behind an authenticating proxy) to allow
# remote admin.
ExecStart=/usr/local/bin/artigate high \
  --listen :8080 \
  --root /var/lib/artigate-high \
  --landing /var/spool/diode-in \
  --public-key /etc/artigate/high.ed25519.pub \
  --import-interval 10s
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/artigate-high /var/spool/diode-in

[Install]
WantedBy=multi-user.target
